# followup_engine/launchd_templates/com.waymark.dailyreport.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.waymark.dailyreport — daily summary email
  Fires at 5:30 PM local time, Monday through Friday. Local time on
  Jamie's Mac is ET, which matches the spec's 5:30 PM ET window.

  TO INSTALL:
    1. cp /Users/jamie/Documents/waymark-hr-group/followup_engine/launchd_templates/com.waymark.dailyreport.plist ~/Library/LaunchAgents/
    2. launchctl load ~/Library/LaunchAgents/com.waymark.dailyreport.plist

  Weekday integers: Sunday = 0, Monday = 1, Tuesday = 2, ..., Saturday = 6.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.waymark.dailyreport</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/jamie/Documents/waymark-hr-group/followup_engine/daily_report.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/jamie/Documents/waymark-hr-group/followup_engine</string>
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Weekday</key><integer>1</integer>
            <key>Hour</key><integer>17</integer>
            <key>Minute</key><integer>30</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>2</integer>
            <key>Hour</key><integer>17</integer>
            <key>Minute</key><integer>30</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>3</integer>
            <key>Hour</key><integer>17</integer>
            <key>Minute</key><integer>30</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>4</integer>
            <key>Hour</key><integer>17</integer>
            <key>Minute</key><integer>30</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>5</integer>
            <key>Hour</key><integer>17</integer>
            <key>Minute</key><integer>30</integer>
        </dict>
    </array>
    <key>StandardOutPath</key>
    <string>/Users/jamie/Documents/waymark-hr-group/followup_engine/logs/daily_report_cron.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jamie/Documents/waymark-hr-group/followup_engine/logs/daily_report_cron.log</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
