# followup_engine/launchd_templates/com.waymark.replydetector.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.waymark.replydetector — reply detection polling loop
  Fires every 2 minutes. Per spec SECTION 9 this gives us sub-5-minute
  push-notification latency on real replies. Pushover keys are loaded
  from followup_engine/.env by the Python script.

  TO INSTALL (after Pushover keys are in .env and tested):
    1. cp /Users/jamie/Documents/waymark-hr-group/followup_engine/launchd_templates/com.waymark.replydetector.plist ~/Library/LaunchAgents/
    2. launchctl load ~/Library/LaunchAgents/com.waymark.replydetector.plist
    3. launchctl list | grep waymark    # verify loaded
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.waymark.replydetector</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/jamie/Documents/waymark-hr-group/followup_engine/reply_detector.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/jamie/Documents/waymark-hr-group/followup_engine</string>
    <key>StartInterval</key>
    <integer>120</integer>
    <key>StandardOutPath</key>
    <string>/Users/jamie/Documents/waymark-hr-group/followup_engine/logs/reply_detector_cron.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jamie/Documents/waymark-hr-group/followup_engine/logs/reply_detector_cron.log</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
